# infra/macos/kg.alistore.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>kg.alistore.backup</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/alistore/bin/alistore-backup.sh</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
		<key>DATABASE_NAME</key>
		<string>alistore_dev</string>
		<key>BACKUP_DIR</key>
		<string>/Users/alistore/backups/alistore</string>
		<key>BACKUP_KEEP_DAYS</key>
		<string>14</string>
	</dict>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>3</integer>
		<key>Minute</key>
		<integer>17</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/alistore/Library/Logs/alistore-backup.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/alistore/Library/Logs/alistore-backup.err.log</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
